# deploy/wmux.service.example
[Unit]
Description=wmux browser terminal multiplexer
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/operator/wmux
Environment=NODE_ENV=development
Environment=WMUX_HOST=100.64.0.10
Environment=WMUX_PORT=3478
Environment=WMUX_CERT_FILE=
Environment=WMUX_KEY_FILE=
Environment=WMUX_PUBLIC_URL=
ExecStartPre=/usr/bin/env node scripts/wmux-wait-bind ${WMUX_HOST}
ExecStart=/usr/bin/npm run dev -- --host ${WMUX_HOST} --port ${WMUX_PORT}
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
