# deploy/wmux.service.example
[Unit]
Description=wmux browser terminal multiplexer
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/gisenberg/git/gisenberg/wmux
Environment=NODE_ENV=development
Environment=WMUX_HOST=100.107.241.79
Environment=WMUX_PORT=3478
ExecStart=/usr/bin/npm run dev -- --host ${WMUX_HOST} --port ${WMUX_PORT}
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
